# conservativeContainerUpdate.service
[Unit]
Description=Conservatively update container images

[Service]
Type=oneshot
ExecStart=%h/containers/conservativeContainerUpdate/conservativeContainerUpdate.py -a immich -f %h/containers/immich/image_versions.env -r immich.target --gotify-url "${GOTIFY_URL}" --gotify-token "${GOTIFY_TOKEN}"
ExecStart=%h/containers/conservativeContainerUpdate/conservativeContainerUpdate.py -a authentik -f %h/containers/authentik/image_versions.env -r authentik.target --gotify-url "${GOTIFY_URL}" --gotify-token "${GOTIFY_TOKEN}"
EnvironmentFile=%h/containers/conservativeContainerUpdate/gotify_creds.env
